# pcr-signature.service
[Unit]
Description=Import TPM2 credentials into the initrd
DefaultDependencies=false

Wants=cryptsetup-pre.target
Before=cryptsetup-pre.target

RequiresMountsFor=/boot/efi

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/pcr-signature.sh

[Install]
WantedBy=cryptsetup.target
